# snapOCR/release/linux/snapOCR.service
[Unit]
Description=SnapOCR Background Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/path/to/start_snapocr.sh
ExecStart=/usr/bin/python3 /path/to/your/snapocr_service.py
Restart=always

[Install]
WantedBy=multi-user.target